# EnviroScripts/rcfcron/crontab
3,18,33,48 * * * * /direct/astro+u/hneal/getnbdump
